# traceability.service
[Unit]
Description=Traceability app (tickets + temperatures)
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=__USER__
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/__USER__/.Xauthority
WorkingDirectory=__APPDIR__
ExecStart=/usr/bin/python3 __APPDIR__/main.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical.target
